# pi-system/wifi-watchdog.timer
[Unit]
Description=Run WiFi watchdog every 5 minutes

[Timer]
OnBootSec=60
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
