# pi-system/wifi-watchdog.service
[Unit]
Description=WiFi watchdog — bounce wlan0 if router unreachable
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi-watchdog.sh
