# evo-system/llama-swap-coder.service
[Unit]
Description=Swap to overnight coder model

[Service]
Type=oneshot
ExecStart=/home/james/llama-swap-coder.sh
