# evo-system/llama-swap-main.service
[Unit]
Description=Swap to daytime main model

[Service]
Type=oneshot
ExecStart=/home/james/llama-swap-main.sh
