# evo-voice/clawdbot-voice.service
[Unit]
Description=Clawd Voice Listener (EVO X2)
After=clawdbot-memory.service
Wants=clawdbot-memory.service

[Service]
Type=simple
User=james
Group=audio
WorkingDirectory=/home/james/clawdbot-memory
ExecStart=/home/james/clawdbot-memory/venv/bin/python /home/james/clawdbot-memory/voice_listener.py
Restart=always
RestartSec=5
Environment=CLAWDBOT_URL=http://10.0.0.1:3000
Environment=EVO_CLASSIFIER_URL=http://127.0.0.1:8081
Environment=DASHBOARD_TOKEN=VhPJmjOLM0A_t2idQrtfa3cHpSr_hBh0fgNxMr2TwUM
Environment=ROUTE_URL=http://localhost:5100/route-command
Environment=MIC_GAIN=3.5
Environment=SPEECH_THRESHOLD=1800
Environment=SILENCE_THRESHOLD=700
Environment=SILENCE_DURATION=0.8
Environment=WHISPER_MODEL=distil-small.en
Environment=LOG_LEVEL=INFO

[Install]
WantedBy=multi-user.target
