# evo-voice/dream-mode.service
[Unit]
Description=Clawd Dream Mode — overnight conversation summarisation
After=network.target
# No hard Requires — dream_mode.py has its own health check wait loop
# that's robust against model swaps, restarts, or any LLM server changes.
Wants=clawdbot-memory.service

[Service]
Type=oneshot
User=james
WorkingDirectory=/home/james/clawdbot-memory
Environment=EVO_LLM_URL=http://localhost:8080
Environment=EVO_MEMORY_URL=http://localhost:5100
Environment=PI_LOG_DIR=/tmp/conversation-logs

# Pull conversation logs from Pi via direct ethernet
ExecStartPre=/bin/bash -c 'mkdir -p /tmp/conversation-logs && scp pi@10.0.0.1:~/clawdbot/data/conversation-logs/*.jsonl /tmp/conversation-logs/ 2>/dev/null || true'

# Run dream mode — script waits for LLM readiness internally (up to 120s)
ExecStart=/home/james/clawdbot-memory/venv/bin/python /home/james/clawdbot-memory/dream_mode.py

StandardOutput=journal
StandardError=journal
TimeoutStartSec=600
